# archive/obsidian-commit.plist
<!--
	This is an example .plist I tried to create for launchd. Ended up running the
	script through crontab instead as I couldn't get this to work well.
-->

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>obsidian-commit</string>
        <key>Program</key>
        <string>/Users/crow/scripts/commit-obsidian.sh</string>
        <key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>10</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
    </dict>
</plist>
